# install/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shell-mcp.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{DAEMON_ENTRY}}</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SHELL_MCP_PORT</key>
    <string>{{PORT}}</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{LOG_DIR}}/daemon.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{{LOG_DIR}}/daemon.stderr.log</string>
</dict>
</plist>
